# SubUpdator/sub-relay/sub-relay/systemd/xray-client.service
[Unit]
Description=Xray Client for Sub-Relay (SOCKS5 Tunnel)
After=network.target
Documentation=https://github.com/XTLS/Xray-core

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/xray run -config /opt/sub-relay/xray-client/config.json
Restart=always
RestartSec=3

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/xray

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xray-client

[Install]
WantedBy=multi-user.target
